# .config/autostart/swapesc.desktop
[Desktop Entry]
Exec=/usr/bin/setxkbmap -option caps:escape
Name=swapesc
Comment=Swap caps lock with escape
Type=Application
